# scripts/vector-mcp.service
[Unit]
Description=Vector MCP Go Backend Daemon
After=network.target

[Service]
Type=simple
User=nilesh
Group=nilesh
WorkingDirectory=/home/nilesh/Documents/vector-mcp/vector-mcp-go
EnvironmentFile=/home/nilesh/Documents/vector-mcp/vector-mcp-go/.env
ExecStart=/home/nilesh/Documents/vector-mcp/vector-mcp-go/vector-mcp-go -daemon
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
